# templates/monasca-forwarder.service.j2
[Unit]
Description=Monasca forwarder
Requires=network-online.target
After=network-online.target monasca-collector

[Service]
Type=simple
Restart=on-failure
User=root
Group={{ monasca_group }}
ExecStart={{ agent_install_dir }}/bin/monasca-forwarder

[Install]
WantedBy=multi-user.target
